# emailwhiz-supervisord.ini
[program:emailwhiz]
command=python manage.py runserver 0.0.0.0:8000
directory=/usr/src/app/EmailWhiz/
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisord/emailwhiz-supervisord.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisord/emailwhiz-supervisord.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5


